# dev_build/systemd/ollama-frontend.service
[Unit]
Description=Ollama Chat Frontend
After=network.target ollama-backend.service
Wants=network.target
Requires=ollama-backend.service

[Service]
Type=simple
User=ec2-user
Group=ec2-user
WorkingDirectory=/home/ec2-user/sample-mpc-app-using-aws-nitrotpm/frontend
Environment=PATH=/usr/local/bin:/usr/bin:/bin
Environment=NODE_ENV=production
Environment=PORT=3000
ExecStart=/usr/bin/node server.js
Restart=always
RestartSec=3
StandardOutput=journal
StandardError=journal

[Install]
WantedBy=multi-user.target